# telegram-forwarder.service
[Unit]
Description=Telegram Forwarder Service
After=network.target

[Service]
Type=simple
User=aka
WorkingDirectory=/home/aka/telegram-forwarder-py
Environment="PATH=/home/aka/telegram-forwarder-py/venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/aka/telegram-forwarder-py/venv/lib/python3.12/site-packages"
Environment="TELEGRAM_API_ID=37956695"
Environment="TELEGRAM_API_HASH=e99150bf9cf5fc5c6ad625a43bae93ed"
ExecStart=/home/aka/telegram-forwarder-py/venv/bin/python /home/aka/telegram-forwarder-py/forwarder.py forward
Restart=on-failure
RestartSec=10

# Security - relaxed for user service
# (Removed restrictions causing GROUP permission errors)

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-forwarder

[Install]
WantedBy=multi-user.target
